# systemd/fixfirst-autopost.service
[Unit]
Description=FixFirst Autopost Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fixfirst-autopost
EnvironmentFile=/opt/fixfirst-autopost/.env
ExecStart=/usr/bin/node /opt/fixfirst-autopost/bridge.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
